# lexacube-nfc-control.service
[Unit]
Description=Lexacube NFC Control Daemon
After=network.target mosquitto.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lexacube
ExecStart=/opt/lexacube/cube_env/bin/python3 /opt/lexacube/nfc_control_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
